# docker/webapp_allinone/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command = uwsgi --emperor /app/files_uwsgi --gid www-data --logto /var/log/uwsgi/emperor.log 

[program:nginx]
command = nginx -g "daemon off;"
